# etc/logrotate/storm-frontend.cron
00 00 * * * root /usr/sbin/logrotate -s /opt/storm/frontend/etc/logrotate/logrotate.status /opt/storm/frontend/etc/logrotate/storm-frontend.logrotate
