# etc/logrotate/storm-frontend.cron
00 00 * * * root /usr/sbin/logrotate -s /opt/srmv2storm/etc/logrotate/logrotate.status /opt/srmv2storm/etc/logrotate/storm-frontend.logrotate
