# etc/cron.d/storm-frontend.cron.in
00 00 * * * root /usr/sbin/logrotate -s @prefix@/etc/logrotate/logrotate.status @prefix@/etc/logrotate/storm-frontend.logrotate
