# etc/cron.d/storm-frontend-server.cron.in
PATH=/sbin:/bin:/usr/sbin:/usr/bin
00 00 * * * root /usr/sbin/logrotate -s @sysconfdir@/logrotate.d/logrotate.status @sysconfdir@/logrotate.d/@PACKAGE@
